# scripts/systemd/airflow-heartbeat.service
[Service]
Type=oneshot
ExecStart=/opt/turbine/airflow-heartbeat.sh
