# scripts/systemd/airflow-confapply-agent.service
[Service]
Type=oneshot
ExecStartPre=/usr/bin/systemctl daemon-reload
ExecStart=/opt/turbine/restart-services.sh
