# scripts/systemd/airflow-heartbeat.timer
[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=airflow.service
