# scripts/systemd/airflow-terminate.service
[Service]
Type=oneshot
ExecStart=/opt/turbine/airflow-terminate.sh
