# scripts/systemd/cfn-hup.service
[Service]
ExecStart=/opt/aws/bin/cfn-hup
Restart=always

[Install]
WantedBy=multi-user.target
